# deploy/budgetpilot.service
[Unit]
Description=BudgetPilot web dashboard
After=network.target

[Service]
Type=simple
User=budgetpilot
Group=budgetpilot
WorkingDirectory=/var/lib/budgetpilot
EnvironmentFile=-/etc/budgetpilot/budgetpilot.env
Environment=BUDGETPILOT_HOME=/var/lib/budgetpilot
Environment=BUDGETPILOT_HOST=127.0.0.1
Environment=BUDGETPILOT_PORT=8765
ExecStart=/opt/budgetpilot/.venv/bin/gunicorn --chdir /opt/budgetpilot --bind ${BUDGETPILOT_HOST}:${BUDGETPILOT_PORT} --workers 2 --timeout 60 budgetpilot_web:app
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/budgetpilot

[Install]
WantedBy=multi-user.target
